# deployments/hello-systemd-extension/extension/etc/systemd/system/hello-world-extension.service
[Unit]
Description=Run a hello-world script provided by a systemd sysext

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/usr/bin/hello-world-extension

[Install]
WantedBy=multi-user.target
