# deployments/hello-service/overlays/etc/systemd/system/hello-world.service
[Unit]
Description=Hello-world service

[Service]
Type=oneshot
ExecStart=echo "hello, world!"

[Install]
WantedBy=multi-user.target
